# packaging/hyprmarker.service
[Unit]
Description=Hyprmarker - Screen annotation tool for Wayland
Documentation=https://github.com/devmobasa/hyprmarker
PartOf=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/bin/hyprmarker --daemon
Restart=on-failure
RestartSec=5
# Ensure Wayland environment is available
Environment="PATH=/usr/local/bin:/usr/bin:/bin"

[Install]
WantedBy=default.target
