# packaging/wayscriber.service
[Unit]
Description=Wayscriber - Screen annotation tool for Wayland
Documentation=https://wayscriber.com
PartOf=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/bin/wayscriber --daemon
Restart=on-failure
RestartSec=5
# Ensure Wayland environment is available
Environment="PATH=/usr/local/bin:/usr/bin:/bin"

[Install]
WantedBy=default.target
